# ops/systemd/iclone-bootstrapper.service
[Unit]
Description=iCLONE ACP Bootstrapper (CLIENT agent)
After=network-online.target iclone-server.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=iclone
Group=iclone
WorkingDirectory=/opt/iclone

ExecStart=/opt/iclone/venv312/bin/python3.12 /opt/iclone/ops/bootstrapper.py \
    --daily-budget 10.0 \
    --max-price 0.10 \
    --jobs-per-hour 10

Restart=on-failure
RestartSec=60

StandardOutput=append:/var/log/iclone/bootstrapper.log
StandardError=append:/var/log/iclone/bootstrapper.log

Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/iclone"
EnvironmentFile=/opt/iclone/.env

LimitNOFILE=65536
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
